# ops/openclaw/systemd/founderos-worker.service.example
[Unit]
Description=Founderos OpenClaw Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/founderos
EnvironmentFile=/root/.config/founderos/aps.env
ExecStart=/usr/bin/env bash /root/.openclaw/workspace/founderos/services/openclaw/worker-loop.sh
Restart=always
RestartSec=5
StandardOutput=append:/var/log/founderos-worker.log
StandardError=append:/var/log/founderos-worker.log

[Install]
WantedBy=multi-user.target
